# public/apps/reference/cv2x-daemon/scripts/cv2x.target
[Unit]
Description=CV2X Applications

[Install]
WantedBy=multi-user.target
